# files/system/usr/etc/xdg/autostart/certux-law.desktop
[Desktop Entry]
Type=Application
Name=Certux Law Configuration
Exec=/usr/bin/certux-law-setup
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
